# dotfiles/.config/autostart/launch.sh.desktop
[Desktop Entry]
Exec=/home/pjones/.config/polybar/launch.sh
Icon=dialog-scripts
Name=launch.sh
Path=
Type=Application
X-KDE-AutostartScript=true
